# crontab/scheduler.sh
#!/bin/bash

date > /last_cron_log.txt

